# docker/crontab
# Patent Law University — Daily Scrapers
# All times UTC.  EDT = UTC-4,  EST = UTC-5.
# 15:00 UTC = 11:00 AM EDT (summer) / 10:00 AM EST (winter)
# Runs Mon–Fri only; skip weekends since courts don't publish on weekends.

0 15 * * 1-5 cd /app && bash /app/scripts/run_daily.sh >> /app/logs/daily.log 2>&1
